# webpage/incavead_make_highscore.service
[Unit]
Description=Generate incavead highscore table

[Service]
Type=oneshot
ExecStart=/srv/http/make_highscore.sh
User=http
