# webpage/incavead_make_highscore.timer
[Unit]
Description=Generate incavead highscore table

[Timer]
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
